# templates/etc/systemd/system/mailman3.service.d/override.conf.j2
[Unit]
After=mysqld.service
[Service]
Restart=on-failure

